# other/remserial.service
[Service]
Type=simple
ExecStart=/usr/bin/remserial -m 2 -d -p 2000 -s "9600 raw" /dev/ttyAMA0 &

[Install]
WantedBy=multi-user.target